# packages/deep-dashboard/deep-dashboard-0.3.1.tar.gz/deep-dashboard-0.3.1/docker/supervisor.conf.d/aiohttp.conf
[program:aiohttp]
numprocs = 4
numprocs_start = 1
process_name = dashboard_%(process_num)s

; Unix socket paths are specified by command line.
command = /usr/local/bin/deep-dashboard --listen-path=/tmp/dashboard_%(process_num)s.sock --config-dir /etc/deep-dashboard/

user = www-data
autostart = true
autorestart = true
#stdout_logfile = /dev/stdout
#stderr_logfile = /dev/stderr
